# build_files/config/systemd/sddm.service.d/gaming-mode.conf
[Unit]
# Override to add gaming-mode target dependency
PartOf=gamemode.target
After=gamemode.target
Conflicts=kodi.target
StopWhenUnneeded=yes
