# build_files/config/systemd/sddm.service.d/gaming-mode.conf
[Unit]
# Override to add gaming-mode target dependency
PartOf=gaming-mode.target
After=gaming-mode.target
Conflicts=kodi-mode.target
StopWhenUnneeded=yes

[Service]
# Ensure SDDM has the correct autologin configuration
ExecStartPre=/usr/bin/configure-sddm-gaming
