# build_files/config/systemd/sddm.service.d/session-modes.conf
[Unit]
PartOf=gaming-mode.target desktop-mode.target
After=gaming-mode.target desktop-mode.target
Conflicts=kodi-mode.target

[Service]
ExecStartPre=/usr/bin/configure-sddm-gaming
