# build_files/config/systemd/gamemode.target
[Unit]
Description=Gaming Mode
Conflicts=kodi.target
After=multi-user.target
Wants=sddm.service
Wants=plugin_loader.service
AllowIsolate=yes

[Install]
WantedBy=graphical.target
