# build_files/config/systemd/session-switch-handler.service
[Unit]
Description=Session Switch Handler
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/session-switch-handler
Restart=always
RestartSec=5

User=root
Group=root

LimitNOFILE=4096
TasksMax=16

PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/run /var/lib/session-state

StandardOutput=journal
StandardError=journal
SyslogIdentifier=session-switch

[Install]
WantedBy=multi-user.target
